# vendor/cookbooks/chattin/files/default/god.conf
# god upstart service

description "God process monitoring"
author "Charles Strahan <charles.c.strahan@gmail.com>"

# When to start and stop  the service
start on runlevel [2345]
stop on runlevel [016]

# Automatically restart process if crashed
respawn

env RACK_ENV=production

# Essentially lets upstart know the process will detach itself to the
# background
expect fork

# Start the process
exec god -l /var/log/god.log -c /usr/local/god/god.conf
